# config/GratiaStaticGraphs.cron
19 9   * * * root /usr/bin/static_graphs >> /var/log/static_graphs.log 2>> /var/log/static_graphs.err
15 */6 * * * root /usr/share/GratiaWeb/gratia_web_auth_update 
