# config/gratiaweb_clean.cron
*/10 * * * * root /usr/bin/gratiaweb_clean_procs
